# etc/systemd/system/getty@tty7.service.d/override.conf
[Service]
Environment=XDG_SESSION_TYPE=x11
