# .config/systemd/user/liskin-timer@.service
[Unit]
Description=liskin timer service template

[Service]
Type=oneshot
ExecStart=%h/bin/liskin-timerwrapper %h/bin/liskin-timer-%i
