# .config/systemd/user/liskin-backup-bash-history.service
[Unit]
Description=liskin-backup-bash-history

[Service]
Type=oneshot
ExecStart=%h/bin/liskin-backup-bash-history
