# .config/systemd/user/liskin-media-mpris-daemon.service
[Unit]
Description=liskin-media mpris-daemon

[Service]
ExecStart=%h/bin/liskin-media mpris-daemon
Restart=always

[Install]
WantedBy=default.target
