# .config/systemd/user/mpd-snapshot.service
[Unit]
Description=snapshot mpd playlist

[Service]
Type=oneshot
ExecStart=%h/bin/mpd-save snapshot
