# .config/systemd/user/rclone-dropbox.service
[Unit]
Description=rclone-dropbox

[Service]
Type=notify
ExecStart=/usr/bin/rclone mount --vfs-cache-mode=full dropbox: %h/mnt/dropbox
