# .config/systemd/user/taskwiki-commit.timer
[Unit]
Description=taskwiki-commit timer

[Timer]
OnCalendar=hourly
RandomizedDelaySec=5m
Persistent=true

[Install]
WantedBy=default.target
