# .config/systemd/user/taskwiki-commit.service
[Unit]
Description=taskwiki-commit

[Service]
Type=oneshot
ExecStart=%h/bin/taskwiki-commit
