# .config/systemd/user/liskin-timer-1d@arbtt-stats-daily.timer.d/override.conf
[Timer]
OnCalendar=
OnCalendar=*-*-* 03:00:00
