# etc/systemd/system/conntrackd.service.d/override.conf
[Unit]
PartOf=nfct-helper-ssdp.service
Requisite=nfct-helper-ssdp.service
After=nfct-helper-ssdp.service

[Install]
WantedBy=nfct-helper-ssdp.service
