# etc/systemd/system/getty@tty9.service.d/override.conf
[Service]
Environment=XDG_SESSION_TYPE=x11 _LISKIN_NVIDIA=1
