# .config/systemd/user/pipewire.service.d/override.conf
[Service]
ExecStartPost=sleep 1
ExecStartPost=%h/bin/pw-unload-module libpipewire-module-x11-bell
