# etc/systemd/system/plocate-updatedb.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/local/sbin/liskin-updatedb-prune
