# .config/systemd/user/dunst.service.d/override.conf
[Service]
Restart=always
ExecStart=
ExecStart=%h/bin/dunst
