# etc/rc.local
#!/usr/bin/env bash

set -eu

/usr/local/sbin/liskin-nvidia off || :
