# etc/profile.d/z-podman-docker.sh
# negate /etc/profile.d/podman-docker.sh
unset DOCKER_HOST
