# etc/systemd/system/anacron.service.d/override.conf
[Service]
# avoid killing forked exim4 trying to deliver results
ExecStart=
ExecStart=/usr/local/sbin/liskin-anacron-wrapper /usr/sbin/anacron -d -q $ANACRON_ARGS
