# etc/systemd/system/anacron.service.d/override.conf
[Service]
# avoid killing forked exim4 trying to deliver results
ExitType=cgroup
