# .config/systemd/user/redshift.service.d/override.conf
[Unit]
StartLimitIntervalSec=0

[Service]
RestartSec=1
ExecStartPre=/bin/sh -c "/usr/bin/xset q >/dev/null 2>&1"
