# .config/systemd/user/blueman-applet.service.d/override.conf
[Unit]
PartOf=graphical-session.target

[Service]
Restart=always
