# .config/systemd/user/liskin-timer@arbtt-rotate.service.d/override.conf
[Service]
Environment=TIMER_WRAPPER_TIMEOUT=600
