# etc/systemd/system/systemd-suspend.service.d/z-liskin-suspend.conf
[Service]
ExecStart=
ExecStart=/usr/local/sbin/liskin-systemd-suspend
