# etc/network/if-up.d/liskin-vpn-route
#!/usr/bin/env bash
ip ro ad $(ip -json ro g 81.31.33.35 | jq -r '.[] | .dst + " via " + .gateway + " dev " + .dev')
