# etc/systemd/system/user@.service.d/override.conf
[Service]
# enable cgroup delegation, needed by https://kind.sigs.k8s.io/docs/user/rootless/
Delegate=yes
